# services-agent/service/autoreply-agent.service
[Unit]
Description=Autoreply Agent Service (tmpfs-root)
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=simple
User=autoreply
Group=autoreply
WorkingDirectory=/opt/autoreply/services-agent
RuntimeDirectory=autoreply-agent
RuntimeDirectoryMode=0755
PIDFile=%t/autoreply-agent/autoreply-agent.pid
ExecStart=/usr/bin/perl /opt/autoreply/services-agent/autoreply-agent.pl
ExecStartPost=/bin/sh -c 'echo $MAINPID > %t/autoreply-agent/autoreply-agent.pid'
ExecStopPost=/bin/sh -c 'rm -f %t/autoreply-agent/autoreply-agent.pid'
EnvironmentFile=/opt/autoreply/env/autoreply-agent.env
UMask=0007

# Grundumgebung
Environment=PATH=/usr/bin:/usr/sbin:/bin:/sbin
UnsetEnvironment=PERL5LIB PERLLIB PERL5OPT IFS CDPATH ENV BASH_ENV LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT

# Restart-Policy
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30

# tmpfs-Root
TemporaryFileSystem=/

# Essentielle Systempfade (schreibgeschützt)
BindReadOnlyPaths=/usr
BindReadOnlyPaths=/lib
BindReadOnlyPaths=/lib64
BindReadOnlyPaths=/bin
BindReadOnlyPaths=/sbin
BindReadOnlyPaths=/etc/resolv.conf
BindReadOnlyPaths=/etc/passwd
BindReadOnlyPaths=/etc/group

# Projektverzeichnisse
BindPaths=/opt/autoreply

# Schreibzugriff für den Service
ReadWritePaths=/opt/autoreply/log      
ReadWritePaths=/opt/autoreply/services-agent/tmp  
ReadWritePaths=/opt/autoreply/script/json
ReadWritePaths=/opt/autoreply/script/conf


# Sandboxing
NoNewPrivileges=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ProtectSystem=strict
ProcSubset=pid
ProtectProc=invisible
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@obsolete
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
KeyringMode=private

# Caps minimal und korrekt
AmbientCapabilities=CAP_CHOWN CAP_FOWNER
CapabilityBoundingSet=CAP_CHOWN CAP_FOWNER


# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autoreply-agent

[Install]
WantedBy=multi-user.target
